# dnsmasq-override.service
[Unit]
Requires=hostapd.service network-online.target
After=hostapd.service network-online.target

[Service]
Restart=on-failure
RestartSec=5s
